# build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

<<<<<<< HEAD
<<<<<<< HEAD
export CATKIN_TEST_RESULTS_DIR="/home/rione/indigo_workspace/project1/build/test_results"
=======
export CATKIN_TEST_RESULTS_DIR="/home/shivarni/indigo_workspace/project1/build/test_results"
>>>>>>> cd51746fa4cc54c811d436b8ef92bbb9eb97c388
=======
export CATKIN_TEST_RESULTS_DIR="/home/rione/indigo_workspace/project1/build/test_results"
>>>>>>> cookingRobot
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
